# roles/elasticsearch/templates/etc/systemd/system/elasticsearch.service.d/override.conf
[Service]
LimitMEMLOCK=infinity
LimitNOFILE=65536
Restart=always
# TODO: swap settings https://www.elastic.co/guide/en/elasticsearch/reference/current/setup-configuration-memory.html
